# roles/entropy/templates/rngd.service.j2
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=/sbin/rngd -f -r {{ rngd_device }}

[Install]
WantedBy=multi-user.target
